# deploy/com.kkabi.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.kkabi</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/유저이름/kkabi/venv/bin/python</string>
        <string>/Users/유저이름/kkabi/main.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/유저이름/kkabi</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/유저이름/kkabi/logs/stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/유저이름/kkabi/logs/stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/Users/유저이름/.npm-global/bin</string>
        <key>GITHUB_TOKEN</key>
        <string>ghp_xxxxx</string>
    </dict>
</dict>
</plist>
